# _staging/etc/systemd/system/shutdown.service
[Unit]
Description=Shutdown and poweroff service
After=mame-autostart.service

[Service]
TTYPath=/dev/tty1
ExecStart=/sbin/poweroff
StandardInput=tty

[Install]
WantedBy=multi-user.target
